# 333home.service
[Unit]
Description=333HOME - Gestionnaire Parc Informatique
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pie333
Group=pie333
WorkingDirectory=/home/pie333/333HOME
Environment=PATH=/home/pie333/333HOME/.venv/bin
ExecStart=/home/pie333/333HOME/.venv/bin/python app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pie333/333HOME

[Install]
WantedBy=multi-user.target